# .github/deploy/silverbank.service
[Unit]
Description=SilverBank Spring Boot Application
After=network.target

[Service]
User=silverbank
Group=silverbank
WorkingDirectory=/home/silverbank/app
ExecStart=/usr/bin/java -jar /home/silverbank/app/silver-0.0.1-SNAPSHOT.jar --spring.profiles.active=prod
SuccessExitStatus=143
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
